# files/etc/rc.local
#!/bin/sh

# 等待网络就绪
sleep 5

# 启动打印机服务
/etc/init.d/cups enable
/etc/init.d/cups start
/etc/init.d/avahi-daemon enable
/etc/init.d/avahi-daemon start

# 设置USB打印机权限
chmod 666 /dev/usb/lp0 2>/dev/null || true
chmod 666 /dev/usb/lp1 2>/dev/null || true

# 设置主机名
uci set system.@system[0].hostname='rk3128-printer'
uci commit system

# 重启网络服务使主机名生效
/etc/init.d/network restart

exit 0
